# deploy/equities/systemd/flux-equities.target
[Unit]
Description=Flux Equities Stack
Wants=flux@equities-api.service
Wants=flux@equities-portfolio.service
Wants=flux@equities-bridge.service
Wants=flux@equities-node-aapl_tradexyz_maker.service
Wants=flux@equities-node-aapl_tradexyz_taker.service
Wants=flux@equities-node-amd_tradexyz_maker.service
Wants=flux@equities-node-amd_tradexyz_taker.service
Wants=flux@equities-node-amzn_tradexyz_maker.service
Wants=flux@equities-node-amzn_tradexyz_taker.service
Wants=flux@equities-node-googl_tradexyz_maker.service
Wants=flux@equities-node-googl_tradexyz_taker.service
Wants=flux@equities-node-meta_tradexyz_maker.service
Wants=flux@equities-node-meta_tradexyz_taker.service
Wants=flux@equities-node-msft_tradexyz_maker.service
Wants=flux@equities-node-msft_tradexyz_taker.service
Wants=flux@equities-node-nvda_tradexyz_maker.service
Wants=flux@equities-node-nvda_tradexyz_taker.service
Wants=flux@equities-node-orcl_tradexyz_maker.service
Wants=flux@equities-node-orcl_tradexyz_taker.service
Wants=flux@equities-node-pltr_tradexyz_maker.service
Wants=flux@equities-node-pltr_tradexyz_taker.service
Wants=flux@equities-node-tsla_tradexyz_maker.service
Wants=flux@equities-node-tsla_tradexyz_taker.service
After=network-online.target

[Install]
WantedBy=multi-user.target
