# deploy/equities/systemd/flux-equities.target
[Unit]
Description=Flux Equities Stack
Wants=flux@equities-api.service
Wants=flux@equities-portfolio.service
Wants=flux@equities-bridge.service
Wants=flux@equities-node-aapl_tradexyz_maker.service
Wants=flux@equities-node-aapl_tradexyz_taker.service
Wants=flux@equities-node-amd_tradexyz_maker.service
Wants=flux@equities-node-amd_tradexyz_taker.service
Wants=flux@equities-node-amzn_binance_perp_maker.service
Wants=flux@equities-node-amzn_binance_perp_taker.service
Wants=flux@equities-node-amzn_tradexyz_maker.service
Wants=flux@equities-node-amzn_tradexyz_taker.service
Wants=flux@equities-node-baba_tradexyz_maker.service
Wants=flux@equities-node-baba_tradexyz_taker.service
Wants=flux@equities-node-coin_binance_perp_maker.service
Wants=flux@equities-node-coin_binance_perp_taker.service
Wants=flux@equities-node-coin_tradexyz_maker.service
Wants=flux@equities-node-coin_tradexyz_taker.service
Wants=flux@equities-node-crcl_binance_perp_maker.service
Wants=flux@equities-node-crcl_binance_perp_taker.service
Wants=flux@equities-node-crcl_tradexyz_maker.service
Wants=flux@equities-node-crcl_tradexyz_taker.service
Wants=flux@equities-node-crwv_tradexyz_maker.service
Wants=flux@equities-node-crwv_tradexyz_taker.service
Wants=flux@equities-node-ewy_binance_perp_maker.service
Wants=flux@equities-node-ewy_binance_perp_taker.service
Wants=flux@equities-node-googl_tradexyz_maker.service
Wants=flux@equities-node-googl_tradexyz_taker.service
Wants=flux@equities-node-hood_binance_perp_maker.service
Wants=flux@equities-node-hood_binance_perp_taker.service
Wants=flux@equities-node-hood_tradexyz_maker.service
Wants=flux@equities-node-hood_tradexyz_taker.service
Wants=flux@equities-node-intc_binance_perp_maker.service
Wants=flux@equities-node-intc_binance_perp_taker.service
Wants=flux@equities-node-intc_tradexyz_maker.service
Wants=flux@equities-node-intc_tradexyz_taker.service
Wants=flux@equities-node-meta_tradexyz_maker.service
Wants=flux@equities-node-meta_tradexyz_taker.service
Wants=flux@equities-node-msft_tradexyz_maker.service
Wants=flux@equities-node-msft_tradexyz_taker.service
Wants=flux@equities-node-mstr_binance_perp_maker.service
Wants=flux@equities-node-mstr_binance_perp_taker.service
Wants=flux@equities-node-mstr_tradexyz_maker.service
Wants=flux@equities-node-mstr_tradexyz_taker.service
Wants=flux@equities-node-mu_tradexyz_maker.service
Wants=flux@equities-node-mu_tradexyz_taker.service
Wants=flux@equities-node-nflx_tradexyz_maker.service
Wants=flux@equities-node-nflx_tradexyz_taker.service
Wants=flux@equities-node-nvda_tradexyz_maker.service
Wants=flux@equities-node-nvda_tradexyz_taker.service
Wants=flux@equities-node-orcl_tradexyz_maker.service
Wants=flux@equities-node-orcl_tradexyz_taker.service
Wants=flux@equities-node-pltr_binance_perp_maker.service
Wants=flux@equities-node-pltr_binance_perp_taker.service
Wants=flux@equities-node-pltr_tradexyz_maker.service
Wants=flux@equities-node-pltr_tradexyz_taker.service
Wants=flux@equities-node-rivn_tradexyz_maker.service
Wants=flux@equities-node-rivn_tradexyz_taker.service
Wants=flux@equities-node-sndk_tradexyz_maker.service
Wants=flux@equities-node-sndk_tradexyz_taker.service
Wants=flux@equities-node-tsla_binance_perp_maker.service
Wants=flux@equities-node-tsla_binance_perp_taker.service
Wants=flux@equities-node-tsla_tradexyz_maker.service
Wants=flux@equities-node-tsla_tradexyz_taker.service
Wants=flux@equities-node-tsm_tradexyz_maker.service
Wants=flux@equities-node-tsm_tradexyz_taker.service
Wants=flux@equities-node-usar_tradexyz_maker.service
Wants=flux@equities-node-usar_tradexyz_taker.service
After=network-online.target

[Install]
WantedBy=multi-user.target
