# deploy/equities/systemd/flux-equities.target
[Unit]
Description=Flux Equities Stack
Wants=flux@equities-api.service
Wants=flux@equities-portfolio.service
Wants=flux@equities-bridge.service
Wants=flux@equities-ibkr-reference-publisher.service
Wants=flux@equities-controller.service
Wants=flux@equities-node-aapl_tradexyz.service
Wants=flux@equities-node-amd_tradexyz.service
Wants=flux@equities-node-amzn_tradexyz.service
Wants=flux@equities-node-googl_tradexyz.service
Wants=flux@equities-node-meta_tradexyz.service
Wants=flux@equities-node-msft_tradexyz.service
Wants=flux@equities-node-nvda_tradexyz.service
Wants=flux@equities-node-orcl_tradexyz.service
Wants=flux@equities-node-pltr_tradexyz.service
Wants=flux@equities-node-tsla_tradexyz.service
Wants=flux@equities-node-amzn_binance_perp.service
Wants=flux@equities-node-coin_binance_perp.service
Wants=flux@equities-node-crcl_binance_perp.service
Wants=flux@equities-node-ewy_binance_perp.service
Wants=flux@equities-node-hood_binance_perp.service
Wants=flux@equities-node-intc_binance_perp.service
Wants=flux@equities-node-mstr_binance_perp.service
Wants=flux@equities-node-pltr_binance_perp.service
Wants=flux@equities-node-tsla_binance_perp.service
After=network-online.target

[Install]
WantedBy=multi-user.target
