# deploy/tokenmm/systemd/flux-tokenmm-telemetry-health.service
[Unit]
Description=TokenMM telemetry health guardrail
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
EnvironmentFile=-/etc/flux/common.env
ExecStart=/bin/bash -lc 'telemetry_dir="${TOKENMM_TELEMETRY_DIR:-/var/lib/nautilus/telemetry/tokenmm}"; cd "${WORKDIR:?}" && exec ./.venv/bin/python ops/scripts/deploy/tokenmm_telemetry_healthcheck.py --telemetry-dir "${telemetry_dir}" --state-db-path "${telemetry_dir}/shipper_state.sqlite" --max-telemetry-dir-gb "${TOKENMM_HEALTH_MAX_DIR_GB:-120}" --max-root-usage-pct "${TOKENMM_HEALTH_MAX_ROOT_PCT:-92}" --max-shipper-lag-minutes "${TOKENMM_HEALTH_MAX_LAG_MIN:-30}"'
StandardOutput=journal
StandardError=journal
SyslogIdentifier=flux-tokenmm-telemetry-health
