# deploy/tokenmm/systemd/flux-tokenmm.target
[Unit]
Description=Flux TokenMM Stack
Wants=flux@tokenmm-api.service
Wants=flux@tokenmm-portfolio.service
Wants=flux@tokenmm-bridge.service
Wants=flux@tokenmm-telemetry-shipper.service
Wants=flux@tokenmm-prometheus.service
Wants=flux@tokenmm-grafana.service
Wants=flux@tokenmm-liquidity-exporter.service
Wants=flux@tokenmm-markouts-exporter.service
Wants=flux@tokenmm-node-plumeusdt_bybit_perp_makerv3.service
Wants=flux@tokenmm-node-plumeusdt_bybit_spot_makerv3.service
Wants=flux@tokenmm-node-plumeusdt_okx_perp_makerv3.service
Wants=flux@tokenmm-node-plumeusdt_binance_perp_makerv3.service
Wants=flux@tokenmm-node-plumeusdt_binance_spot_makerv3.service
Wants=flux@tokenmm-node-plumeusdt_bitget_perp_makerv3.service
Wants=flux@tokenmm-node-plumeusdt_bitget_spot_makerv3.service
After=network-online.target

[Install]
WantedBy=multi-user.target
